# system/horus-kbd-boot.service
[Unit]
Description=Horus: color del teclado (tema) antes del login
# El hidraw puede tardar en frio (misma carrera EC que horus-bat-limit);
# el wrapper hace su propio poll de 30s, no dependemos de udev aqui.

[Service]
Type=oneshot
Environment=HORUS_USER_HOME=@HOME@
ExecStart=/usr/local/bin/horus-kbd-boot

[Install]
WantedBy=multi-user.target
