# system/horus-gpu-watch.service
[Unit]
Description=Horus Project — GPU de sesion segun perfil de energia
[Service]
ExecStart=%h/.local/bin/horus-gpu-watch
Restart=on-failure
[Install]
WantedBy=default.target
